# systemd/amprgate-poll.service
[Unit]
Description=AMPRNet Gateway Status Poller
Documentation=https://github.com/KE4CON/FieldCommand-IMS
After=network.target fcc-lookup.service
Wants=network.target

[Service]
Type=simple
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/amprgate_poll.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amprgate-poll

[Install]
WantedBy=multi-user.target
