# systemd/deadmans.service
[Unit]
Description=FieldCommand Dead Man's Switch Monitor
After=fcc-lookup.service
Requires=fcc-lookup.service

[Service]
Type=simple
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/deadmans.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deadmans

[Install]
WantedBy=multi-user.target
