# systemd/fcc-lookup.service
[Unit]
Description=FieldCommand FCC Lookup Server (Port 5050)
After=network.target
Wants=network.target

[Service]
Type=simple
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/fcc_lookup_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fcc-lookup

[Install]
WantedBy=multi-user.target
