# systemd/fcc-refresh.service
[Unit]
Description=FieldCommand FCC Database Refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/build_fcc_db.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fcc-refresh
