# systemd/fieldcommand-backup@.service
[Unit]
Description=FieldCommand USB Backup — %i
After=media-fieldcommand.mount
Requires=media-fieldcommand.mount

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash -c 'DEST=/media/fieldcommand/backup/$(date +%%Y%%m%%d_%%H%%M%%S) && mkdir -p "$DEST/data" "$DEST/incidents" && rsync -a /opt/fieldcommand/data/ "$DEST/data/" && sqlite3 /opt/fieldcommand/data/fieldcommand.db ".backup $DEST/data/fieldcommand_backup.db" && echo "Backup complete: $DEST"'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fieldcommand-backup
