# systemd/health-monitor.service
[Unit]
Description=FieldCommand Health Monitor (Port 5051)
After=network.target
Wants=network.target

[Service]
Type=simple
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/health_monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-monitor

[Install]
WantedBy=multi-user.target
