# systemd/ics-platform.service
[Unit]
Description=FieldCommand ICS Platform Server (Port 5055)
After=network.target
Wants=network.target

[Service]
Type=simple
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/ics_platform_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ics-platform

[Install]
WantedBy=multi-user.target
