# systemd/pat.service
[Unit]
Description=Pat Winlink — Backup Winlink Client (Port 8090)
Documentation=https://getpat.io
After=network.target
Wants=network.target

[Service]
Type=simple
# Pat runs as the fieldcommand user so it can access its config at
# /home/fieldcommand/.config/pat/config.json
User=fieldcommand
Group=fieldcommand
ExecStart=/usr/bin/pat http
Restart=on-failure
RestartSec=10
# Give Pat a clean environment
Environment=HOME=/opt/fieldcommand
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pat

[Install]
WantedBy=multi-user.target
