# systemd/repeater-refresh.service
[Unit]
Description=FieldCommand Repeater Database Refresh (RepeaterBook)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=fieldcommand
Group=fieldcommand
WorkingDirectory=/opt/fieldcommand
# REPEATERBOOK_TOKEN is read from the environment file below.
# Create /opt/fieldcommand/data/repeaterbook.env with a line:
#   REPEATERBOOK_TOKEN=your-issued-token
EnvironmentFile=-/opt/fieldcommand/data/repeaterbook.env
ExecStart=/opt/fieldcommand/venv/bin/python /opt/fieldcommand/python/fetch_repeaters.py --states IL,WI,IN,IA --bands 2m,70cm
StandardOutput=journal
StandardError=journal
SyslogIdentifier=repeater-refresh
